# roles/systemd/templates/unit-undefined.service.j2
{{ ansible_managed | comment }}
[Unit]
Description={{ item.name }}

[Service]
ExecStart={{ systemd_script_install_path }}/{{ item.name }}

[Install]
WantedBy=multi-user.target
