# roles/laptop/templates/root-suspend.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Local system suspend actions
Before=sleep.target

[Service]
{{ root_suspend }}

[Install]
WantedBy=sleep.target
