# files/gitlab/docker-auto-prune.service
[Unit]
Description=Prune Docker
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker system prune -af
