# files/gitlab/podman-auto-prune.service
[Unit]
Description=Prune Podman for gitlab-runner

[Service]
User=gitlab-runner
Group=gitlab-runner
ExecStart=/usr/bin/podman system prune -af
