# files/gitlab/registry-garbage-collect.service
[Unit]
Description=Gitlab garbage collection
After=gitlab-runsvdir.service

[Service]
ExecStart=/usr/bin/gitlab-ctl registry-garbage-collect
