# roles/netclient/templates/netclient.timer.j2
[Unit]
Description=Calls the Netmaker Mesh Client Service
Requires=netclient.service

[Timer]
Unit=netclient.service
OnCalendar={{ netclient_timer_interval }}

[Install]
WantedBy=timers.target
