# roles/borg/templates/backup.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Borg backup for {{ item.name }}

[Service]
ExecStart={{ borg_backup_script_path }}/backup-{{ item.name }}

[Install]
WantedBy=multi-user.target
