# files/wireguard/wireguard_reresolve-dns.timer.j2
[Unit]
Description=Periodically reresolve DNS of all WireGuard endpoints

[Timer]
OnCalendar={{ wiregaurd_reresolve_dns_oncalendar | default('hourly') }}

[Install]
WantedBy=timers.target
