# roles/withings_sync/files/withings-sync.timer
[Unit]
Description=Sync Withings to Garmin every 4 hours

[Timer]
Persistent=yes
OnCalendar=*-*-* 00/4:15:00

[Install]
WantedBy=timers.target
